# bootstrap/launchd/com.homelab.colima.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent: starts the Colima container runtime when the user logs in.
  Installed by bootstrap/05-enable-autostart.sh into ~/Library/LaunchAgents/.

  Logs: /tmp/homelab-colima.{log,err}
  Paths assume Intel Mac with Homebrew at /usr/local/bin.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homelab.colima</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/colima</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/homelab-colima.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/homelab-colima.err</string>
</dict>
</plist>
