# bootstrap/launchd/com.homelab.k3d.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent: starts the k3d 'homelab' cluster after Colima is reachable,
  then runs bootstrap/heal.sh to recover from the well-known Colima-restart
  wedge where the k3d serverlb crash-loops on
      nginx: [emerg] host not found in upstream "k3d-homelab-agent-X:443"
  heal.sh detects that and falls back to a set of alpine/socat bridges so
  the cluster is reachable on its expected host ports (6445 for API,
  80/443/53 for ingress + DNS).

  Logs: /tmp/homelab-k3d.{log,err}
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homelab.k3d</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>for i in {1..60}; do /usr/local/bin/docker info >/dev/null 2>&amp;1 &amp;&amp; break; sleep 2; done; /usr/local/bin/k3d cluster start homelab; sleep 5; HEAL=&quot;$HOME/src/homelab-iac/bootstrap/heal.sh&quot;; [ -x &quot;$HEAL&quot; ] &amp;&amp; &quot;$HEAL&quot; || true</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/homelab-k3d.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/homelab-k3d.err</string>
</dict>
</plist>
