# ansible/roles/host/templates/com.homelab.dns-proxy.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Managed by Ansible. Do not hand-edit. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.homelab.dns-proxy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/usr/local/lib/homelab/dns-proxy.py</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>DNS_LISTEN</key>
    <string>{{ lan_ip }}</string>
    <key>DNS_UPSTREAM</key>
    <string>{{ colima_vm_ip }}</string>
    <key>DNS_PORT</key>
    <string>53</string>
  </dict>

  <key>UserName</key>
  <string>root</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/var/log/homelab-dns-proxy.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/homelab-dns-proxy.log</string>
</dict>
</plist>
