# bootstrap/launchd-system/com.homelab.mac-extras.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchDaemon: mac-extras.sh, the textfile-collector script that gathers
  macOS-specific metrics node_exporter's darwin build doesn't expose
  (CPU temp, fan RPM, SMART, colima status, launchd health, ...).

  Installed by bootstrap/10-mac-exporter.sh into /Library/LaunchDaemons.

  Why root: powermetrics, smartctl, and ioreg need it. Running as root
  also means the produced files are owned by root, matching the
  node_exporter LaunchDaemon that reads them.

  Why StartInterval 60: prometheus scrape interval is 30s, but a script
  that shells out to powermetrics + smartctl per disk is not free. 60s
  is enough fidelity for thermal/wear trends. Bump if you ever need it.

  COLIMA_USER is stamped at install time by 10-mac-exporter.sh so the
  script can `sudo -u <user> colima status` without prompting (colima
  config lives in the user's home dir, not root's).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.homelab.mac-extras</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/usr/local/lib/homelab/mac-extras.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin</string>
    <key>OUT_DIR</key>
    <string>/usr/local/var/mac-exporter</string>
    <key>COLIMA_USER</key>
    <string>__COLIMA_USER__</string>
    <key>COLIMA_PROFILE</key>
    <string>__COLIMA_PROFILE__</string>
  </dict>

  <key>UserName</key>
  <string>root</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>60</integer>

  <!-- One-shot per tick; don't restart on exit. -->
  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/var/log/homelab-mac-extras.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/homelab-mac-extras.log</string>
</dict>
</plist>
