# bootstrap/launchd-system/com.homelab.node-exporter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchDaemon: prometheus node_exporter (darwin build, from Homebrew).
  Installed by bootstrap/10-mac-exporter.sh into /Library/LaunchDaemons.

  Why a LaunchDaemon (system) and not a LaunchAgent (user-session):
    - Survives logout. The Mac is a headless server; metrics must stay up
      whether or not anyone is logged in.
    - Runs as root, which is what the textfile collector files written
      by com.homelab.mac-extras (also root) expect. Avoids any "Operation
      not permitted" on reads when we tighten file modes later.

  Why bind to __BIND_ADDR__:
    bootstrap/10-mac-exporter.sh discovers the Mac's vmnet bridge IP
    (the address `host.docker.internal` resolves to from inside Colima)
    and stamps it here at install time. That address is private to the
    Colima VM — not reachable from the home LAN — so we get a "loopback-
    like" exposure without losing reachability from in-cluster pods.

  Why expose all collectors except hwmon/textfile defaults:
    The darwin build only ships a small set anyway (cpu, diskstats,
    filesystem, loadavg, meminfo, netdev, textfile, time, uname).
    Explicit --collector.* flags keep the surface predictable when the
    formula bumps and a new default collector appears.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.homelab.node-exporter</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE_EXPORTER__</string>
    <string>--web.listen-address=__BIND_ADDR__:9100</string>
    <string>--collector.textfile.directory=/usr/local/var/mac-exporter</string>
    <!-- The darwin build of node_exporter ships only a small set of
         collectors (cpu, diskstats, filesystem, loadavg, meminfo, netdev,
         textfile, time, uname). Don't pass --no-collector.* flags here:
         most linux-only collectors aren't compiled in, and the binary
         rejects unknown flags outright. -->
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>UserName</key>
  <string>root</string>

  <!-- Don't start at boot — Colima must be up first for the vmnet bind
       address to exist. start.sh kicks this daemon after Colima is ready. -->
  <key>RunAtLoad</key>
  <false/>

  <!-- Restart on crash, but with backoff so a broken update doesn't
       hammer the system. -->
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/var/log/homelab-node-exporter.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/homelab-node-exporter.log</string>
</dict>
</plist>
