# bootstrap/launchd-system/com.homelab.cluster.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchDaemon: starts the full homelab stack at boot (Colima → k3d → heal).
  Installed by bootstrap/05-enable-autostart.sh into /Library/LaunchDaemons/.
  Runs as homelab-admin (Colima/Docker need the user's $HOME).

  Logs: /tmp/homelab-cluster.{log,err}
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homelab.cluster</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>/Users/homelab-admin/src/homelab-iac/bootstrap/start.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/homelab-admin</string>
    </dict>

    <key>UserName</key>
    <string>homelab-admin</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/homelab-cluster.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/homelab-cluster.err</string>
</dict>
</plist>
